# tags/REL-1_23_4-Win/DEBIAN_build/prerm
#!/bin/bash
update-rc.d -f hsflowd remove
update-rc.d -f sflowovsd remove
