# trunk/DEBIAN_build/prerm
#!/bin/bash
update-rc.d hsflowd remove
update-rc.d sflowovsd remove
