# DEBIAN_build/postinst
#!/bin/bash
update-rc.d hsflowd defaults
if [ -e /usr/bin/ovs-vsctl ]; then update-rc.d sflowovsd defaults; fi 
