# DEBIAN_build/postinst
#!/bin/bash
# either:
# update-rc.d hsflowd defaults
# or:
# systemctl daemon-reload
# systemctl enable hsflowd
exit 0
