# tags/REL-1_25_5/DEBIAN_build/prerm
#!/bin/bash
update-rc.d -f hsflowd remove
test -e /etc/init.d/sflowovsd && update-rc.d -f sflowovsd remove
exit 0
